# wallpapers/.config/systemd/user/wallpaper.timer
[Unit]
Description=Periodically change wallpaper

[Timer]
OnBootSec=15s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
