# x/.config/systemd/user/i3lock.service
[Unit]
Description=Lock screen before suspend

[Service]
Type=forking
Environment=DISPLAY=:0
ExecStart=%h/lock.sh
ExecStartPost=/bin/sleep 1
