# config/systemd/user/xinit.target
[Unit]
Description=Xinitr replacement
Wants=xorg.target
