# jira-dump/.config/systemd/user/jira-dump.timer
[Unit]
Description = period dump from jira

[Install]
WantedBy = timers.target

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1800
Persistent=true
